# .config/systemd/user/hyprshade.timer
[Unit]
Description=Apply screen filter on schedule

[Timer]


[Install]
WantedBy=timers.target
